# crontabs/root
# Run every 2 hour
# For more help - https://crontab.guru/
0 */2 * * * /system/bin/detachyt

