# crontabs/root

# Run every 1 hour
0 */1 * * * /system/bin/detachyt

