# crontabs/root

# Run every 6 hours
0 */6 * * * /system/bin/detachyt

